# tools/rpc-gateway/ops/receipt-watch@.service
[Unit]
Description=RPC receipt watcher (%i)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# Template file: replace placeholders before production use.
# Required placeholders:
# - __WORKDIR__
# - __RUN_USER__
# - __RUN_GROUP__
# - __RPC_URL__
# - __WATCH_SCRIPT__
WorkingDirectory=__WORKDIR__
User=__RUN_USER__
Group=__RUN_GROUP__
Environment=EVM_RPC_URL=__RPC_URL__
# Environment=ALERT_WEBHOOK_URL=https://example.com/webhook
ExecStart=__WATCH_SCRIPT__ %i 180 1500
NoNewPrivileges=true
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
