# ops/systemd/kasane-indexer.service.d/10-alert.conf
[Unit]
OnFailure=kasane-alert@%n
OnFailureJobMode=replace-irreversibly
